# debian/dpkg.cron.daily
#!/bin/sh

/usr/libexec/dpkg/dpkg-db-backup
